# supervisor/memcached.conf
[program:memcached]
user=memcache
command=/usr/bin/memcached
autorestart=true
